# crontab
@reboot /home/jc/projects/bigj.dev/deploy.sh
